# ubuntu/crontab
0 0 * * * /home/william/.linuxbrew/bin/zsh /home/william/dotfiles/ubuntu/update.zsh > /tmp/stdout.log  2>/tmp/stderr.log
