# macOS/crontab
0 0 * * * zsh /Users/William/dotfiles/cronmacupdate.zsh > /tmp/stdout.log 2>/tmp/stderr.log
